# .config/systemd/user/cache-bk-timer.service
[Unit]
Description=Timer for backing up tmpfs mounted $USER/.cache

[Timer]
OnUnitActiveSec=900
OnBootSec=900

[Install]
WantedBy=timers.target
