# user.crontab
# minute	hour	day_of_month	month	day_of_week	command
*/5		*	*		*	*		tmpfs-cache-bk
*/15		*	*		*	*		curl -o /tmp/forecast wttr.in/$(cat /tmp/zip)
*		*/4	*		*	*		aur-update-check > /tmp/aur.updates
*		*/4	*		*	*		checkupdates > /tmp/off.updates
0		0	*		*	*		web-connect-extras
