# .config/root.crontab
# minute	hour	day_of_month	month	day_of_week	command
0		0	1		*	*		/home/austin/git/scripts/shell/maintenance
0		11	*		*	1		tar -czf /home/austin/mega/backup-pc/system/etc.tar.gz /etc
0		11	*		*	1		tar -czf /home/austin/mega/backup-pc/system/boot.tar.gz /boot

